# overlay/etc/systemd/system/pommed.service
[Unit]
Description=Apple keyboard support

[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=/usr/sbin/pommed -f
