# overlay/etc/systemd/system/tor.service
[Unit]
Description=tor

[Service]
User=user
ExecStart=/usr/sbin/tor
